# scripts/sambil-frontend.service
[Unit]
Description=Sambil frontend (Next.js standalone)
After=network.target

[Service]
Type=simple
User=git
Group=git
WorkingDirectory=/home/git/sambil
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3010
ExecStart=/usr/bin/node /home/git/sambil/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
